# timers/pacman_update.timer
[Unit]
Description=Run updates daily

[Timer]
OnCalendar=12:00
Persistent=true

[Install]
WantedBy=network-online.target
