# timers/poweroff.timer
[Unit]
Description=Shutdown on friday nights.

[Timer]
OnCalendar=friday 18:00

[Install]
WantedBy=timers.target
